# com.ccgl.growlink.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ccgl.growlink.hourly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dangillan/Documents/Claude/Projects/CCGL Growlink Analyst/push_live.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dangillan/Documents/Claude/Projects/CCGL Growlink Analyst</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dangillan/Documents/Claude/Projects/CCGL Growlink Analyst/logs/growlink-hourly.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dangillan/Documents/Claude/Projects/CCGL Growlink Analyst/logs/growlink-hourly.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
